# deter-report-env/dailyharvesting.cron
# m h  dom mon dow   command
50 08 * * * /usr/local/cron_exec.sh