# deter-report-env/dailyharvesting.cron
# m h  dom mon dow   command
30 07 * * * /usr/local/cron_exec.sh