# raster-process-env/process-schedule.cron
# m h  dom mon dow   command
10,30,50 * * * * /usr/local/cron_exec.sh