# cerrado-deter-env/dailyharvesting.cron
# m h  dom mon dow   command
30 23 * * * /usr/local/cron_exec.sh