# deter-ftp-download-env/ftp-schedule.cron
# m h  dom mon dow   command
59 06 * * * /usr/local/cron_exec.sh