# raster-publisher-env/publisher-schedule.cron
# m h  dom mon dow   command
5,25,45 * * * * /usr/local/cron_exec.sh