# raster-publisher-env/publisher-schedule.cron
# m h  dom mon dow   command
10 * * * * /usr/local/cron_exec.sh