# raster-ftp-download-env/ftp-schedule.cron
# m h  dom mon dow   command
0,20,40 * * * * /usr/local/cron_exec.sh