# raster-ftp-download-env/ftp-schedule.cron
# m h  dom mon dow   command
30 * * * * /usr/local/cron_exec.sh